# deploy/hearth-agents.service
[Unit]
Description=Hearth Agents Autonomous Development Pipeline
After=network.target

[Service]
Type=simple
User=hearth
Group=hearth
WorkingDirectory=/opt/hearth-agents
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=10

Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=512

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hearth-agents

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hearth-agents/logs /opt/hearth-agents/data

[Install]
WantedBy=multi-user.target
